# dns/rootfs/etc/services.d/bind/run
#!/usr/bin/with-contenv bashio
# shellcheck disable=SC1008
# ==============================================================================
# Start bind
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

if [ ! -d /data/cache ]; then
    bashio::log.info "Creating cache folder"
    mkdir -p /data/cache
    chown named:named /data/cache
    chmod 755 /data/cache
fi

## Run your program
exec /usr/sbin/named -g -c /config/named.conf -u named
