# dcsys/rootfs/etc/services.d/dcsysfluviusd/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Declare variables
declare message
declare ha_hostname
declare dcsys_config_root_dir
declare dcsysfluvuisd_args

## Get the user config options.
fluvius_hostname=$(bashio::config 'fluvius_hostname')
dcsysfluviusd_args=$(bashio::config 'dcsysfluvuisd_args')
dcsys_config_root_dir=$(bashio::config 'dcsys_config_root_dir')
bashio::log.info "${fluvius_hostname:="fluvius.struyve.lan"}"

#exec /usr/bin/my_program
export DCSYSROOT=/etc/dcsys
bashio::log.info "$DCSYSROOT/bin/dcsysfluviusd -h ${fluvius_hostname} ${dcsysfluviusd_args}"
$DCSYSROOT/bin/dcsysfluviusd -h ${fluvius_hostname} ${dcsysfluviusd_args}
